# mariadb.ini
[program:mariadb]
priority=1
startsecs=20
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --user=mysql
user=root
autostart=true
autorestart=unexpected
stopsignal=QUIT
stderr_logfile=/var/log/mariadb_err.log
